# .local/var/sv/syncthing/finish
#!/bin/sh
exec syncthing cli operations shutdown
