# .local/var/sv/redshift/run
#!/bin/sh
read -r lat lon < ~/.local/private/location
echo "$lat:$lon"
exec redshift -l "$lat:$lon"
