# .local/var/sv/redshift/run
#!/bin/sh
read -r lat lon _ < ~/.local/private/location
echo "$lat:$lon"
exec redshift -l "$lat:$lon"
